# roles/kraken_systemd/templates/kraken-layercake.service.j2
# container-{{ container_id }}.service
# autogenerated by Podman 3.1.0-dev
# Tue Apr 20 20:19:17 UTC 2021

[Unit]
Description=Podman container-{{ "kraken-layercake-virt" if kr_with_vbox else "kraken-layercake" }}service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start -a {{ container_id }}
ExecStop=/usr/bin/podman stop -t 10 {{ container_id }}
ExecStopPost=/usr/bin/podman stop -t 10 {{ container_id }}
PIDFile=/run/containers/storage/overlay-containers/{{ container_id }}/userdata/conmon.pid

Environment=NOTIFY_SOCKET=/run/systemd/notify
NotifyAccess=all
Type=notify

[Install]
WantedBy=multi-user.target default.target